# nginx-cache-dir.service
[Unit]
Description=Create Nginx FastCGI Cache Directory in RAM
After=local-fs.target
Before=nginx.service

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'mkdir -p /dev/shm/nginx-cache && chown www-data:www-data /dev/shm/nginx-cache && chmod 755 /dev/shm/nginx-cache'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
